# openrc/suwayomi-server
#!/sbin/openrc-run
# /etc/init.d/suwayomi-server
description="A free and open source manga reader server that runs extensions built for Tachiyomi."

user="suwayomi"
group="suwayomi"
root_dir="/var/lib/suwayomi-server"

java_bin="/usr/lib64/openjdk-25/bin/java"
java_args="-Dsuwayomi.tachidesk.config.server.initialOpenInBrowserEnabled=false -Dsuwayomi.tachidesk.config.server.systemTrayEnabled=false -Dsuwayomi.tachidesk.config.server.rootDir=${root_dir} --enable-native-access=ALL-UNNAMED"
jar="/opt/suwayomi-server/suwayomi-server.jar"

export LD_LIBRARY_PATH="/usr/lib64/openjdk-25/lib/server:/usr/lib64/openjdk-25/lib${LD_LIBRARY_PATH:+:${LD_LIBRARY_PATH}}"
export LD_PRELOAD="${root_dir}/bin/kcef/libcef.so"

command="$java_bin"
command_args="$java_args -jar $jar"
command_user="${user}:${group}"
command_background=true

pidfile="/var/run/${RC_SVCNAME}.pid"
output_logger="logger -t suwayomi"
error_logger="logger -t suwayomi -s"

depend() {
    need net
    use logger
}

start_pre() {
    checkpath --directory --owner "${user}:${group}" --mode 0750 "${root_dir}"
}
